# systemd-files/picamera-monitor.service
[Unit]
Description=Picamera Streamer Health Check

[Service]
Type=oneshot
ExecStart=/home/lee/picamera-streamer/health_check.sh
